# crontab
*/5 * * * * python /$PROCESS_NAME.py >> /logs/logfile.log
