# crontab
*/5 * * * * python /vaccine_valley.py >> /logs/logfile.log
  